# deploy/cch-pull-server.service.example
[Unit]
Description=CCH Local Pull WebSocket Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/path/to/cch-local-pull
ExecStart=/usr/bin/env node /path/to/cch-local-pull/src/server.js --config /path/to/cch-local-pull/config/server.json
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
